# src/crontab
* */1 * * * . /src/loadenv.sh && python /src/scripts/generate.py >> /var/log/cron.log 2>&1
* 0 * * * . /src/loadenv.sh && python /src/scripts/train/tfidf.py >> /var/log/cron.log 2>&1
