# src/crontab
* */1 * * * . /src/loadenv.sh && python -m generate >> /var/log/cron.log 2>&1
* 0 * * * . /src/loadenv.sh && python -m tfidf >> /var/log/cron.log 2>&1
0 0 * * 0 . /src/loadenv.sh && python -m distilbert >> /var/log/cron.log 2>&1
