# src/crontab
* */1 * * * /bin/bash -c '. /src/loadenv.sh && python -m generate >> /var/log/cron/generate/$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1'
* 0 * * * /bin/bash -c '. /src/loadenv.sh && python -m tfidf >> /var/log/cron/tfidf/$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1'
0 0 * * 0 /bin/bash -c '. /src/loadenv.sh && python -m distilbert >> /var/log/cron/distilbert/$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1'
